# services/netatalk/run
#!/bin/bash

# wait for avahi-daemon
sv start avahi || exit 1

# run in foreground
exec netatalk -d
